# nav_board.service
[Unit]
Description=Differential GPS MAVLink Service
After=network.target

[Service]
Type=simple
User=pi

# systemd MUST start in the folder where the Pipfile lives
WorkingDirectory=/home/pi/Differential_GPS

# Use the absolute path to pipenv (from Step 6.1), followed by "run python"
ExecStart=/usr/bin/pipenv run python nav_board.py

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target